# producer/crontab
PATH=/usr/local/bin:/usr/bin:/bin
* * * * * python3 -u /app/producer.py >> /proc/1/fd/1 2>> /proc/1/fd/2
